# roles/osbackupper/templates/openstack-vm-backup.timer.j2
[Unit]
Description=Daily OpenStack VM backup timer

[Timer]
OnCalendar=*-*-* {{ osbackupper.run_at }}
Persistent=true

[Install]
WantedBy=timers.target
